# examples/systemd/datronis-relay.service
[Unit]
Description=datronis-relay — chat bridge for the Claude Agent SDK
Documentation=https://github.com/mhmdevan/datronis-relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=datronis
Group=datronis
WorkingDirectory=/opt/datronis-relay

# HOME points inside ReadWritePaths so `claude login` (run once, manually,
# before enabling the service) can write OAuth credentials to
# /var/lib/datronis-relay/home/.claude. ProtectHome=true stays enabled
# because /var/lib is outside the /home, /root, /run/user paths it protects.
Environment=HOME=/var/lib/datronis-relay/home
Environment=DATRONIS_CONFIG_PATH=/etc/datronis-relay/config.yaml

# Optional env file — lets operators inject DATRONIS_TELEGRAM_BOT_TOKEN,
# DATRONIS_SLACK_*, and (optionally) ANTHROPIC_API_KEY without editing the
# unit. The leading `-` makes the file optional; the service boots fine
# without it if subscription login is the only auth in use.
EnvironmentFile=-/etc/datronis-relay/datronis-relay.env

ExecStart=/opt/datronis-relay/venv/bin/datronis-relay
Restart=on-failure
RestartSec=5s

# --- Hardening ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/datronis-relay
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
